# actividad4/hola.service
[Unit]
Description=Hola y Fecha

[Service]
Type=simple
ExecStart=/usr/bin/bash /home/hola.sh

[Install]
WantedBy=multi-user.target